# data/ubuntu-location-provider-geoclue2.conf
description "Ubuntu Geoclue2 Position Provider"

start on started dbus and started ofono and started network-manager
stop on stopping dbus

respawn

exec /usr/bin/ubuntu-location-provider-geoclue2
